# sdk_oem_1.22.0/package/base-files/files/etc/init.d/tracert.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

SERVICE_USE_PID=1

START=51

start() {
	insmod tracert
}

stop() {
	rmmod tracert
}


